# roles/configure-back/templates/back.service.j2
[Unit]
Description=Signup back
After=network-online.target

[Service]
Type=simple

User={{ app_user }}
Group={{ app_user }}

WorkingDirectory=/opt/apps/signup.api.gouv.fr-back/current
EnvironmentFile=/etc/signup.api.gouv.fr-back.conf
ExecStart=/usr/local/bin/bundle exec puma -C ./config/puma.rb

Restart=on-failure
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
